# deploy/model050426-export-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.model050426.export-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/jhbvdnsbkvnsd/Library/Application Support/MODEL050426-sync/bin/sync_exports.sh</string>
  </array>
  <key>RunAtLoad</key>
  <false/>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>EnvironmentVariables</key>
  <dict>
    <key>MODEL050426_REPORT_SCRIPT</key>
    <string>/Users/jhbvdnsbkvnsd/Library/Application Support/MODEL050426-sync/report.py</string>
    <key>MODEL050426_LOCAL_PYTHON</key>
    <string>/Library/Frameworks/Python.framework/Versions/3.11/bin/python3</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/jhbvdnsbkvnsd/Library/Logs/model050426-export-sync.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/jhbvdnsbkvnsd/Library/Logs/model050426-export-sync.err</string>
</dict>
</plist>
