# deploy/model050426.service
[Unit]
Description=MODEL050426 Signal Engine
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/MODEL05042026
EnvironmentFile=/opt/MODEL05042026/.env
ExecStart=/opt/MODEL05042026/.venv/bin/python /opt/MODEL05042026/main.py
Restart=always
RestartSec=5
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
